# scripts/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-telegram-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/claude-telegram-bot/node_modules/.bin/tsx</string>
        <string>/path/to/claude-telegram-bot/src/index.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/claude-telegram-bot</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/claude-telegram-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-telegram-bot.log</string>
</dict>
</plist>
